# main/supervisord.conf
[program:node]
command = node app/server/bootstrap.js
directory = /home/dotcloud/current