# main/supervisord.conf
[program:node]
command = node app/server/server.js
directory = /home/dotcloud/current
